# edge/iams-edge.service
[Unit]
Description=IAMS Edge Device (Face Detection)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/iams-edge
Environment="PATH=/home/pi/iams-edge/venv/bin"
ExecStart=/home/pi/iams-edge/venv/bin/python /home/pi/iams-edge/run.py

# Restart policy (aggressive for edge device reliability)
Restart=always
RestartSec=5
StartLimitInterval=120
StartLimitBurst=10

# Resource limits (conservative for Raspberry Pi)
MemoryLimit=512M
CPUQuota=80%

# Camera access
SupplementaryGroups=video

# Logging
StandardOutput=append:/home/pi/iams-edge/logs/service.log
StandardError=append:/home/pi/iams-edge/logs/service-error.log

[Install]
WantedBy=multi-user.target
